# rootfs/etc/services.d/flood/run
#!/bin/sh
echo "starting flood"
cd /srv/flood/app
npm start
