# rootfs/etc/services.d/flood/run
#!/usr/bin/with-contenv bash

echo "starting flood"
cd /srv/flood/app
s6-setuidgid mediadepot npm start
