# rootfs/etc/services.d/flood/run
#!/usr/bin/with-contenv bash

echo "starting flood"
su -c "cd /srv/flood/app && npm start" mediadepot
