# root/etc/services.d/rutorrent/finish
#!/usr/bin/execlineb -S0

echo "rtorrent exiting"
s6-svscanctl -t /var/run/s6/services