# root/etc/services.d/flood/run
#!/bin/sh
echo "starting flood"
cd /srv/rtorrent/app/flood
npm start
