# root/etc/services.d/flood/run
#!/bin/sh
echo "starting flood"
cd /usr/flood
npm start
